# deploy/smartshell.service.example
# SmartShell — example systemd unit (copy and edit paths, then enable).
#
#   sudo cp deploy/smartshell.service.example /etc/systemd/system/smartshell.service
#   sudo nano /etc/systemd/system/smartshell.service   # fix User, paths
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now smartshell.service
#   sudo systemctl status smartshell.service
#
# Logs: journalctl -u smartshell.service -f

[Unit]
Description=SmartShell accident alert (Phase 1)
After=network-online.target smartshell-buzzer-silence.service
Wants=network-online.target
# Uncomment the next line once smartshell-buzzer-silence.service is enabled:
# Requires=smartshell-buzzer-silence.service

[Service]
Type=simple
User=smartshell
Group=smartshell
SupplementaryGroups=gpio dialout
WorkingDirectory=/home/smartshell/AccidentAlertSystem
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/smartshell/AccidentAlertSystem/.venv/bin/python -m src.main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
